# services/biketag-renew-assets-cert.timer
[Unit]
Description=Run BikeTag Daily Cert Renewal

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
