# services/biketag-db-backup.service
[Unit]
Description=BikeTag DB Backup

[Service]
TimeoutStartSec=0
ExecStart=/home/core/biketag-ops/bin/archive-db 

